# src/go/tmpl/templates/miniccc/miniccc.service
[Unit]
Description=miniccc Agent

[Service]
ExecStart=/usr/local/bin/miniccc -serial /dev/virtio-ports/cc

[Install]
WantedBy=multi-user.target
